# lib/systemd/system/fan.service
# KHr:
# Type=idle  ->  Start after all services have started

[Unit]
Description=fan control service

[Service]
Type=idle
ExecStart=/opt/nastia-server/bin/fan
User=pi

[Install]
WantedBy=multi-user.target

