# lib/systemd/system/ups.service
# KHr:
# Type=idle  ->  Start after all services have started

[Unit]
Description=ups control service

[Service]
Type=idle
ExecStart=/opt/nastia-server/sbin/ups
User=root

[Install]
WantedBy=multi-user.target

