# etc/cron.d/nastia-server
#
# /etc/cron.d/nastia-server
#
# crontab fragment for the nastia-server
# KHr
#

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Disable automated emails
MAILTO="" 


 
# Jobs executed by root
#       m   h   dom mon dow   user   command
        0 */3   *   *   *     root   /opt/nastia-server/sbin/backup-config.sh > /dev/null
        0   1   *   *   *     root   /opt/nastia-server/sbin/permissions.sh   > /dev/null
        5   1   *   *   Sun   root   /opt/nastia-server/sbin/backup-sd.sh     > /dev/null

# Jobs executed by USER
#       m   h   dom mon dow   user  command
      */5   *   *   *   *     pi    /opt/nastia-server/bin/dyndns.sh          > /dev/null
        5 */3   *   *   *     pi    /opt/nastia-server/bin/backup-hdd.sh      > /dev/null
  5-59/10   *   *   *   *     pi    /opt/nastia-server/bin/dropbox-photos.sh  > /dev/null
     */10   *   *   *   *     pi    /opt/nastia-server/bin/photostream.sh     > /dev/null
        0   2   *   *   Sun   pi    /opt/nastia-server/bin/check-images.sh    > /dev/null
        0   4   *   *   *     pi    /opt/nastia-server/bin/monitor.sh         > /dev/null
 
# Test job
#     */1   *   *   *   *     pi    echo "$(date)  shell:$0  path:$PATH  mailto:$MAILTO" >> /home/pi/cron-test

